# processes.ini
[program:netdisco-daemon-fg]
command=localenv netdisco-daemon-fg
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopsignal=TERM

[program:netdisco-web-fg]
command=localenv starman netdisco-web-fg --workers=2 --disable-keepalive --access-log=/dev/stdout --error-log=/dev/stdout
directory=%(ENV_NETDISCO_HOME)s/perl5/bin
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
